# service/openclaw-mini.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.mini</string>

    <key>ProgramArguments</key>
    <array>
      <string>__PROJECT_DIR__/.venv/bin/openclaw-mini</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/service/openclaw-mini.out.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/service/openclaw-mini.err.log</string>
  </dict>
</plist>
